# config/launchagents/io.maridb.r2sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>io.maridb.r2sync</string>

  <key>ProgramArguments</key>
  <array>
    <string>REPLACE_WITH_UV_BIN</string>
    <string>run</string>
    <string>--project</string>
    <string>REPLACE_WITH_PROJECT_DIR</string>
    <string>python</string>
    <string>REPLACE_WITH_PROJECT_DIR/scripts/sync_r2.py</string>
    <string>push-ais-parquet</string>
    <string>--regions</string>
    <string>japansea,blacksea,middleeast,singapore,europe,gulfofmexico,hornofafrica,persiangulf,gulfofaden,gulfofguinea</string>
    <string>--data-dir</string>
    <string>REPLACE_WITH_HOME/.maridb/data/raw/ais</string>
    <string>--staging-dir</string>
    <string>REPLACE_WITH_HOME/.maridb/data/staging/ais</string>
  </array>

  <key>WorkingDirectory</key>
  <string>REPLACE_WITH_PROJECT_DIR</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>AWS_ACCESS_KEY_ID</key>
    <string>REPLACE_WITH_AWS_ACCESS_KEY_ID</string>
    <key>AWS_SECRET_ACCESS_KEY</key>
    <string>REPLACE_WITH_AWS_SECRET_ACCESS_KEY</string>
    <key>S3_BUCKET</key>
    <string>maridb-public</string>
    <key>HOME</key>
    <string>REPLACE_WITH_HOME</string>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:REPLACE_WITH_HOME/.local/bin:REPLACE_WITH_HOME/.cargo/bin</string>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
  </dict>

  <!-- Run every hour; skips DBs whose mtime hasn't changed -->
  <key>StartInterval</key>
  <integer>3600</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>REPLACE_WITH_LOG_DIR/r2sync.log</string>
  <key>StandardErrorPath</key>
  <string>REPLACE_WITH_LOG_DIR/r2sync.err</string>
</dict>
</plist>
