# config/launchagents/io.maridb.aisstream.gulfofaden.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>io.maridb.aisstream.gulfofaden</string>

  <key>ProgramArguments</key>
  <array>
    <string>REPLACE_WITH_UV_BIN</string>
    <string>run</string>
    <string>--project</string>
    <string>REPLACE_WITH_PROJECT_DIR</string>
    <string>python</string>
    <string>-m</string>
    <string>pipelines.ingest.ais_stream</string>
    <string>--db</string>
    <string>REPLACE_WITH_PROJECT_DIR/data/raw/ais/gulfofaden.duckdb</string>
    <string>--bbox</string>
        <string>10</string>
    <string>42</string>
    <string>16</string>
    <string>52</string>
    <string>--duration</string>
    <string>0</string>
    <string>--flush-interval</string>
    <string>60</string>
  </array>

  <key>WorkingDirectory</key>
  <string>REPLACE_WITH_PROJECT_DIR</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>AISSTREAM_API_KEY</key>
    <string>REPLACE_WITH_AISSTREAM_API_KEY</string>
    <key>DB_PATH</key>
    <string>REPLACE_WITH_PROJECT_DIR/data/raw/ais/gulfofaden.duckdb</string>
    <key>HOME</key>
    <string>REPLACE_WITH_HOME</string>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:REPLACE_WITH_HOME/.local/bin:REPLACE_WITH_HOME/.cargo/bin</string>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
  </dict>

  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>30</integer>
  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>REPLACE_WITH_LOG_DIR/gulfofaden.log</string>
  <key>StandardErrorPath</key>
  <string>REPLACE_WITH_LOG_DIR/gulfofaden.err</string>
</dict>
</plist>
